# conf.d/supervisord.conf
[supervisord]  
nodaemon=true  
  
[program:redis]
command=/etc/init.d/redis-server start
  
[program:sd]
command=php /home/SwooleDistributed/bin/start_swoole_server.php start -d